# tests/network/image/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:curl]
command=/bin/sh -c "while true; do sleep 5; curl -m 5 `cat /var/www/html/user-data.html` -o /var/www/html/received-data.html; done"

[program:ip]
command=/bin/sh -c "while true; do sleep 5; ifconfig>/var/www/html/ifconfig.html; done"

[program:dhcpcd]
command=/sbin/dhcpcd --nobackground -f /dhcpcd.conf